# .config/crontab/file
# user crontab

# man 1 crontab for acceptable formats:
#    <minute> <hour> <day> <month> <dow> <tags and command>
#    <@freq> <tags and command>

# reload normal newsboat two times a day
0 */12 * * * /usr/local/bin/newsup normal

# reload other newsboat once every hour
0 * * * * /usr/local/bin/newsup other

# end
