# .local/sv/wlsunset/run
#!/bin/sh
# Turin coordinates
exec wlsunset -l 45.1 -L 7.7
# Lausanne coordinates
# exec wlsunset -l 46.5 -L 6.6
